# deploy/systemd/stocktrading-db-backup.timer
[Unit]
Description=Nightly stock trading PostgreSQL backup

[Timer]
OnCalendar=*-*-* 02:30:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
