# deploy/systemd/stocktrading-job-report.timer
[Unit]
Description=Daily report and dashboard (10:07 ET, market days)

[Timer]
OnCalendar=Mon..Fri *-*-* 10:07:00
Timezone=America/New_York
Persistent=true
Unit=stocktrading-job@report.service

[Install]
WantedBy=timers.target
